# pkgbuilds/base-config/nut-driver-override.service
[Unit]
StopWhenUnneeded=no

[Service]
PIDFile=/run/nut/usbhid-ups-cyberpower.pid
Group=nut
User=nut
RestartSec=10s