# systemd/prometheus-backup.timer
[Unit]
Description=Backup prometheus data every 15min to RAID array and external storage

[Timer]
OnBootSec=15min

[Install]
WantedBy=timers.target