# pkgbuilds/docker-service/container.slice
[Unit]
Description=Container and other similar technology slices
Before=slices.target


[Slice]
AllowedCPUs=0-15