# pkgbuilds/scryped-service/opt-scrypted-large.automount
[Unit]
Description=Automatic mounting of surveillance folder when accessed

[Automount]
Where=/opt/scrypted/large
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target