# pkgbuilds/scryped-service/opt-backups.automount
[Unit]
Description=Automount backups

[Automount]
Where=/opt/backups
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target