# pkgbuilds/scryped-service/post.install
post_install() {
    systemctl enable --now opt-backups.automount
    systemctl enable --now opt-scrypted-large.automount
}