# root/etc/systemd/system/backup.timer
[Unit]
Description=Full System Backup to ssd

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target