# root/etc/systemd/system/ffmpeg.service
[Unit]
Description=Backs up files

[Service]
Type=notify
ExecStart=/data/test/UnixSocket/RicardoAlcaraz.Services.UdpStreamListener