# home/.config/systemd/user/stream-webcam.socket
[Unit]
Description=Socket for xorg at display %i

[Socket]
ListenStream=/